# rootfs/etc/crontabs/apache
*/15 * * * * /usr/bin/php -f /app/cron.php > /dev/null 2>&1
